# lab4/my.cron
1 * * * * /home/dnc28/cs480/lab4/peekaboo.c > /home/dnc28/cs480/lab4/saveaboo.txt
