# lab4/my.cron
* * * * * cd /home/dnc28/cs480/lab4 && ./peekaboo > ./saveaboo.txt
